# rootfs/etc/services.d/xinetd/run
#!/usr/bin/with-contenv bash
exec /usr/sbin/xinetd -dontfork -stayalive -filelog /var/log/xinetd.log