# rootfs/etc/s6-overlay/s6-rc.d/svc-xinetd/run
#!/usr/bin/with-contenv sh

exec "${XINETD_PATH:?}" -dontfork -stayalive -filelog /var/log/xinetd.log -pidfile /var/run/xinetd.pid
